# etc/systemd/system/snapshot-daily.timer
[Unit]
Description=Daily snapshot timer

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
